# src/s6-overlay/s6-rc.d/dbus/run
#!/usr/bin/env bash

# Set up logging
exec 2>&1

mkdir -p /run/dbus

# Start dbus in foreground mode
exec dbus-daemon --system --nofork --nopidfile
