# src/s6-overlay/s6-rc.d/confd/up
#!/command/with-contenv sh
# shellcheck shell=sh

exec /etc/s6-overlay/scripts/configure-balena.sh
