# src/s6-overlay/s6-rc.d/avahi-daemon/run
#!/usr/bin/env bash

# Avahi daemon service for mDNS/DNS-SD
# Runs avahi-daemon in foreground mode for s6 supervision

set -euo pipefail

# Redirect all future stdout/stderr to s6-log
exec > >(exec s6-log -b p"avahi-daemon[$$]:" 1 || true) 2>&1

# Wait for dbus
count=0
until [[ -S /run/dbus/system_bus_socket ]]; do
    sleep 1
    count=$((count + 1))
    if [[ ${count} -gt 5 ]]; then
        echo "Failed to find dbus socket"
        exit 1
    fi
done

# Start avahi-daemon in foreground
exec avahi-daemon --no-chroot
