# src/certs-watch.path
[Path]
PathModified=/certs

[Install]
WantedBy=basic.target
